# tasks/pass/files/config/profile.d/pass.sh
# pass.sh
PASSWORD_STORE_DIR="${XDG_CONFIG_HOME}/password-store"
