# tasks/vcpkg/files/config/profile.d/vcpkg.sh
# vcpkg.sh
export VCPKG_ROOT="${XDG_DATA_HOME}/vcpkg"
